# test/files/etc/uci-defaults/00-shutdown
echo force reboot that should shutdown qemu >> /dev/kmsg

reboot -fn
